# build/packaging/agent/systemd/serviceradar-agent.service
[Unit]
Description=ServiceRadar Agent Service
After=network.target

[Service]
Type=simple
User=serviceradar
WorkingDirectory=/var/lib/serviceradar
ExecStart=/usr/local/bin/serviceradar-agent
Environment="CONFIG_SOURCE=file"
Environment="KV_ADDRESS=127.0.0.1:50057"
Environment="KV_SEC_MODE=mtls"
Environment="KV_CERT_DIR=/etc/serviceradar/certs"
Environment="KV_CERT_FILE=/etc/serviceradar/certs/component.pem"
Environment="KV_KEY_FILE=/etc/serviceradar/certs/component-key.pem"
Environment="KV_CA_FILE=/etc/serviceradar/certs/ca-chain.pem"
EnvironmentFile=-/etc/serviceradar/kv-overrides.env
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
Restart=always
RestartSec=10
LimitNPROC=512
LimitNOFILE=65535

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
