# build/packaging/core-elx/systemd/serviceradar-core-elx.service
[Unit]
Description=ServiceRadar Core Elixir service (core-elx)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/serviceradar/core-elx.env
WorkingDirectory=/usr/local/share/serviceradar-core-elx
ExecStart=/usr/local/share/serviceradar-core-elx/bin/serviceradar_core_elx start
ExecStop=/usr/local/share/serviceradar-core-elx/bin/serviceradar_core_elx stop
Restart=on-failure
User=root
Group=root
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
