# build/packaging/agent-gateway/systemd/serviceradar-agent-gateway.service
[Unit]
Description=ServiceRadar Agent Gateway service (agent-gateway)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/serviceradar/agent-gateway.env
WorkingDirectory=/usr/local/share/serviceradar-agent-gateway
ExecStart=/usr/local/share/serviceradar-agent-gateway/bin/serviceradar_agent_gateway start
ExecStop=/usr/local/share/serviceradar-agent-gateway/bin/serviceradar_agent_gateway stop
Restart=on-failure
User=root
Group=root
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
