# packaging/dusk-checker/systemd/serviceradar-dusk-checker.service
[Unit]
Description=ServiceRadar Dusk Checker Service
After=network.target

[Service]
Type=simple
User=serviceradar
ExecStart=/usr/local/bin/serviceradar-dusk-checker -config /etc/serviceradar/checkers/dusk.json
Restart=always
RestartSec=10
LimitNPROC=512
LimitNOFILE=65535

StandardOutput=append:/var/log/serviceradar/dusk-checker.log
StandardError=append:/var/log/serviceradar/dusk-checker-error.log

[Install]
WantedBy=multi-user.target
