# packaging/kv/systemd/serviceradar-kv.service
[Unit]
Description=ServiceRadar KV Service
After=network.target

[Service]
Type=simple
User=serviceradar
ExecStart=/usr/local/bin/serviceradar-kv -config /etc/serviceradar/kv.json
Restart=always
RestartSec=10
LimitNPROC=512
LimitNOFILE=65535

StandardOutput=append:/var/log/serviceradar/kv.log
StandardError=append:/var/log/serviceradar/kv-error.log

[Install]
WantedBy=multi-user.target
