# packaging/otel/systemd/serviceradar-otel.service
[Unit]
Description=ServiceRadar OTEL Collector
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/serviceradar-otel --config /etc/serviceradar/otel.toml
Restart=on-failure
User=serviceradar
Group=serviceradar

[Install]
WantedBy=multi-user.target