# packaging/sysmonosx_host/scripts/postinstall
#!/bin/bash
# Copyright 2024 Carver Automation Corporation.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

# Post-install script for ServiceRadar sysmon-osx macOS package
# This script is executed after the package files are installed.

set -e

LOG_DIR="/var/log/serviceradar"
PLIST="/Library/LaunchDaemons/com.serviceradar.sysmonosx.plist"
SERVICE_LABEL="com.serviceradar.sysmonosx"
SERVICE_TARGET="system/${SERVICE_LABEL}"
CONFIG="/usr/local/etc/serviceradar/sysmon-osx.json"
CONFIG_BACKUP="/usr/local/etc/serviceradar/.sysmon-osx.json.pre-upgrade"

log() {
    echo "[postinstall] $*"
}

log "Starting sysmon-osx service setup..."

# Create log directory if it doesn't exist
if [ ! -d "${LOG_DIR}" ]; then
    log "Creating log directory: ${LOG_DIR}"
    mkdir -p "${LOG_DIR}"
fi

# Ensure proper permissions on log directory
chmod 755 "${LOG_DIR}"

# Restore pre-upgrade config if it exists (preserves user customizations)
if [ -f "${CONFIG_BACKUP}" ]; then
    log "Restoring previous config from backup"
    mv -f "${CONFIG_BACKUP}" "${CONFIG}"
else
    log "No previous config found, using package default"
fi

# Verify plist exists
if [ ! -f "${PLIST}" ]; then
    log "ERROR: LaunchDaemon plist not found at ${PLIST}"
    exit 1
fi

# Verify binary exists and is executable
BINARY="/usr/local/libexec/serviceradar/serviceradar-sysmon-osx"
if [ ! -x "${BINARY}" ]; then
    log "ERROR: Binary not found or not executable at ${BINARY}"
    exit 1
fi

# Verify config exists
if [ ! -f "${CONFIG}" ]; then
    log "ERROR: Config file not found at ${CONFIG}"
    exit 1
fi

# Stop and unload existing service if running (ignore errors)
log "Stopping existing service (if any)..."
launchctl bootout "${SERVICE_TARGET}" 2>/dev/null || true

# Small delay to ensure clean unload
sleep 1

# Load the new service
log "Loading service from ${PLIST}..."
launchctl bootstrap system "${PLIST}"

# Enable the service to start at boot
log "Enabling service for auto-start..."
launchctl enable "${SERVICE_TARGET}"

# Start the service immediately
log "Starting service..."
launchctl kickstart -k "${SERVICE_TARGET}"

# Verify service is running
sleep 2
if launchctl print "${SERVICE_TARGET}" >/dev/null 2>&1; then
    log "Service started successfully"
else
    log "WARNING: Service may not have started. Check logs at ${LOG_DIR}/sysmon-osx.log"
fi

log "sysmon-osx service setup complete"
exit 0
