# elixir/serviceradar_poller/rel/serviceradar-poller.service
[Unit]
Description=ServiceRadar Poller (Elixir)
Documentation=https://github.com/carverauto/serviceradar
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
User=serviceradar
Group=serviceradar

# Working directory and paths
WorkingDirectory=/opt/serviceradar/poller
ExecStart=/opt/serviceradar/poller/bin/serviceradar_poller start
ExecStop=/opt/serviceradar/poller/bin/serviceradar_poller stop
ExecReload=/opt/serviceradar/poller/bin/serviceradar_poller restart

# Environment configuration
EnvironmentFile=-/etc/serviceradar/poller.env

# Default environment variables
Environment=RELEASE_DISTRIBUTION=name
Environment=CLUSTER_ENABLED=true
Environment=CLUSTER_STRATEGY=epmd
Environment=POLLER_PARTITION_ID=default
Environment=SPIFFE_CERT_DIR=/etc/serviceradar/certs
Environment=ERL_CRASH_DUMP=/var/log/serviceradar/poller_crash.dump

# Restart configuration
Restart=on-failure
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Allow read access to cert directory
ReadOnlyPaths=/etc/serviceradar/certs

# Allow write access to log directory
ReadWritePaths=/var/log/serviceradar

# Capabilities for ICMP ping (if needed)
AmbientCapabilities=CAP_NET_RAW

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=serviceradar-poller

[Install]
WantedBy=multi-user.target
