# packaging/device-mgr/systemd/serviceradar-device-mgr.service
[Unit]
Description=ServiceRadar Device Manager
After=network.target

[Service]
Type=simple
User=serviceradar
ExecStart=/usr/local/bin/serviceradar-device-mgr
Restart=always
RestartSec=10
LimitNOFILE=65535
LimitNPROC=65535

[Install]
WantedBy=multi-user.target