# packaging/profiler/systemd/serviceradar-profiler.service
[Unit]
Description=ServiceRadar Profiler agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/serviceradar-profiler --config /etc/serviceradar/profiler.toml
Environment="CONFIG_SOURCE=kv"
Environment="KV_ADDRESS=127.0.0.1:50057"
Environment="KV_SEC_MODE=mtls"
Environment="KV_CERT_DIR=/etc/serviceradar/certs"
Environment="KV_CERT_FILE=/etc/serviceradar/certs/profiler.pem"
Environment="KV_KEY_FILE=/etc/serviceradar/certs/profiler-key.pem"
Environment="KV_CA_FILE=/etc/serviceradar/certs/root.pem"
EnvironmentFile=-/etc/serviceradar/kv-overrides.env
Restart=on-failure
User=serviceradar
Group=serviceradar

StandardOutput=append:/var/log/serviceradar/profiler.log
StandardError=append:/var/log/serviceradar/profiler-error.log

[Install]
WantedBy=multi-user.target
