# packaging/core/systemd/serviceradar-core.service
[Unit]
Description=ServiceRadar Core Service
After=network.target

[Service]
Type=simple
User=serviceradar
ExecStart=/usr/local/bin/serviceradar-core -config /etc/serviceradar/core.json
EnvironmentFile=/etc/serviceradar/api.env
Environment="CONFIG_SOURCE=kv"
Environment="KV_ADDRESS=127.0.0.1:50057"
Environment="KV_SEC_MODE=mtls"
Environment="KV_CERT_DIR=/etc/serviceradar/certs"
Environment="KV_CERT_FILE=/etc/serviceradar/certs/core.pem"
Environment="KV_KEY_FILE=/etc/serviceradar/certs/core-key.pem"
Environment="KV_CA_FILE=/etc/serviceradar/certs/root.pem"
EnvironmentFile=-/etc/serviceradar/kv-overrides.env
Restart=always
RestartSec=10
TimeoutStopSec=20
KillMode=mixed
KillSignal=SIGTERM
StandardOutput=append:/var/log/serviceradar/core.log
StandardError=append:/var/log/serviceradar/core-error.log

[Install]
WantedBy=multi-user.target
