# elixir/serviceradar_agent/rel/serviceradar-agent.service
[Unit]
Description=ServiceRadar Agent (Elixir)
Documentation=https://github.com/carverauto/serviceradar
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
User=serviceradar
Group=serviceradar

# Working directory and paths
WorkingDirectory=/opt/serviceradar/agent
ExecStart=/opt/serviceradar/agent/bin/serviceradar_agent start
ExecStop=/opt/serviceradar/agent/bin/serviceradar_agent stop
ExecReload=/opt/serviceradar/agent/bin/serviceradar_agent restart

# Environment configuration
EnvironmentFile=-/etc/serviceradar/agent.env

# Default environment variables
Environment=RELEASE_DISTRIBUTION=name
Environment=CLUSTER_ENABLED=true
Environment=CLUSTER_STRATEGY=epmd
Environment=AGENT_PARTITION_ID=default
Environment=SPIFFE_CERT_DIR=/etc/serviceradar/certs
Environment=ERL_CRASH_DUMP=/var/log/serviceradar/agent_crash.dump

# Restart configuration
Restart=on-failure
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Allow read access to cert directory
ReadOnlyPaths=/etc/serviceradar/certs

# Allow write access to log directory
ReadWritePaths=/var/log/serviceradar

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=serviceradar-agent

[Install]
WantedBy=multi-user.target
