# packaging/rperf-checker/systemd/serviceradar-rperf-checker.service
[Unit]
Description=ServiceRadar RPerf Network Performance Checker
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/serviceradar-rperf-checker --config /etc/serviceradar/checkers/rperf.json
Environment="CONFIG_SOURCE=kv"
Environment="KV_ADDRESS=127.0.0.1:50057"
Environment="KV_SEC_MODE=mtls"
Environment="KV_CERT_DIR=/etc/serviceradar/certs"
Environment="KV_CERT_FILE=/etc/serviceradar/certs/rperf-checker.pem"
Environment="KV_KEY_FILE=/etc/serviceradar/certs/rperf-checker-key.pem"
Environment="KV_CA_FILE=/etc/serviceradar/certs/root.pem"
EnvironmentFile=-/etc/serviceradar/kv-overrides.env
Restart=on-failure
User=serviceradar
Group=serviceradar
StandardOutput=append:/var/log/serviceradar/rperf-checker.log
StandardError=append:/var/log/serviceradar/rperf-checker-error.log

[Install]
WantedBy=multi-user.target
