# packaging/rperf-checker/systemd/serviceradar-rperf-checker.service
[Unit]
Description=ServiceRadar RPerf Network Performance Checker
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/serviceradar-rperf-checker --config /etc/serviceradar/checkers/rperf.json
Restart=on-failure
User=serviceradar
Group=serviceradar
StandardOutput=append:/var/log/serviceradar/rperf-checker.log
StandardError=append:/var/log/serviceradar/rperf-checker-error.log

[Install]
WantedBy=multi-user.target