# packaging/web-ng/systemd/serviceradar-web-ng.service
[Unit]
Description=ServiceRadar Phoenix web (web-ng)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/serviceradar/web-ng.env
WorkingDirectory=/usr/local/share/serviceradar-web-ng
ExecStart=/usr/local/share/serviceradar-web-ng/bin/serviceradar_web_ng start
ExecStop=/usr/local/share/serviceradar-web-ng/bin/serviceradar_web_ng stop
Restart=on-failure
User=root
Group=root
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
